# .local/etc/cron.daily/dotfiles
#!/usr/bin/env sh

echo "Updating ${HOME}/Dotfiles."
test -x /usr/bin/git && git -C ${HOME}/Dotfiles pull
